# pi-client/mesophy-media-daemon.service
[Unit]
Description=Mesophy Digital Signage Media Daemon with Advanced Display Management
After=network.target graphical-session.target systemd-modules-load.service
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=pi
Group=video
WorkingDirectory=/opt/mesophy/pi-client
ExecStart=/usr/bin/node /opt/mesophy/pi-client/media-daemon.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=NODE_ENV=production
Environment=DISPLAY=:0
Environment=FRAMEBUFFER=/dev/fb0
Environment=LD_PRELOAD=/usr/lib/arm-linux-gnueabihf/libbrcmEGL.so:/usr/lib/arm-linux-gnueabihf/libbrcmGLESv2.so

# Resource limits
LimitNOFILE=65536
MemoryMax=512M

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/mesophy /boot /boot/firmware /dev/fb0 /sys/class/drm
PrivateTmp=true

# Required for framebuffer, media access, and display management
SupplementaryGroups=audio video dialout i2c spi gpio

[Install]
WantedBy=multi-user.target